# home/Library/LaunchAgents/com.user.mackup-backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.user.mackup-backup</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/mackup</string>
    <string>backup</string>
    <string>-f</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>12</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/tmp/mackup-backup.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/mackup-backup.err</string>
</dict>
</plist>
